# deploy/listpull.service
[Unit]
Description=ListPull Decklist Application
Documentation=https://github.com/your-org/listpull
After=network.target

[Service]
Type=simple
User=listpull
Group=listpull
WorkingDirectory=/opt/listpull
ExecStart=/usr/bin/node server.js
Restart=on-failure
RestartSec=10

# Environment
EnvironmentFile=/opt/listpull/.env

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/listpull

# Resource limits
MemoryMax=512M
CPUQuota=100%

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=listpull

[Install]
WantedBy=multi-user.target
